# autostart/openvpn-tray.desktop
[Desktop Entry]
Type=Application
Name=OpenVPN Tray
Comment=VPN Manager System Tray Application
Exec=%h/.local/bin/openvpn-tray
Icon=network-vpn
Terminal=false
Categories=Network;System;
StartupNotify=false
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel
