# systemd/openvpn-tray.service
[Unit]
Description=OpenVPN System Tray Application
Documentation=https://github.com/loujr/archtools
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/openvpn-tray
Restart=on-failure
RestartSec=5

# Environment for X11/Wayland display access
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=graphical-session.target
