# systemd/shroud.service
# Shroud VPN Manager - Systemd User Service
#
# This service file runs the Shroud VPN tray application as a user service.
# It requires a graphical session with D-Bus access.
#
# Installation:
#   1. Copy this file to ~/.config/systemd/user/shroud.service
#   2. Run: systemctl --user daemon-reload
#   3. Run: systemctl --user enable --now shroud.service
#
# Verification:
#   - Check status: systemctl --user status shroud.service
#   - View logs: journalctl --user -u shroud.service -f
#
# Notes:
#   - This service requires the graphical-session.target to be active
#   - For Wayland/KDE Plasma, DISPLAY and WAYLAND_DISPLAY are typically inherited
#   - If the tray icon doesn't appear, ensure your DE supports StatusNotifierItem
#
# Uninstallation:
#   1. Run: systemctl --user disable --now shroud.service
#   2. Remove: rm ~/.config/systemd/user/shroud.service
#   3. Run: systemctl --user daemon-reload

[Unit]
Description=Shroud VPN Manager
Documentation=https://github.com/loujr/archtools
After=graphical-session.target
Wants=graphical-session.target
# Ensure D-Bus session is available
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
# Path to the binary - adjust if installed elsewhere
ExecStart=%h/.local/bin/openvpn-tray
# Restart on failure, but not on clean exit (e.g., user quit)
Restart=on-failure
RestartSec=5
# Maximum restart attempts (5 within 60 seconds)
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment for graphical session
# Most DEs export these automatically, but we can be explicit
# Uncomment if the tray icon doesn't appear:
# Environment=DISPLAY=:0
# Environment=WAYLAND_DISPLAY=wayland-0

# XDG runtime directory (usually set by login manager)
# Environment=XDG_RUNTIME_DIR=/run/user/%U

# Logging (optional - logs go to journald by default)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shroud

[Install]
WantedBy=graphical-session.target
